# systemd/traffic-guard-update.service
[Unit]
Description=Update Traffic Guard ipset blacklist
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/update-traffic-guard.sh
